# edushare.service
[Unit]
Description=EduShare Gunicorn Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/edushare
Environment="PATH=/var/www/edushare/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=edushare_project.settings_production"
EnvironmentFile=/var/www/edushare/.env.production
ExecStart=/var/www/edushare/venv/bin/gunicorn \
    --config /var/www/edushare/gunicorn_config.py \
    edushare_project.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
