# edushare.service
[Unit]
Description=EduShare Gunicorn Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/home/react-edu-share-production
Environment="PATH=/home/react-edu-share-production/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=edushare_project.settings_production"
EnvironmentFile=/home/react-edu-share-production/.env.production
ExecStart=/home/react-edu-share-production/venv/bin/gunicorn \
    --config /home/react-edu-share-production/gunicorn_config.py \
    edushare_project.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
